# overlay/etc/supervisor/conf.d/log-watcher.conf
[program:log-watcher]
command=/scripts/log-watcher.sh
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
